# system_files/base/usr/lib/systemd/user/orora-cli.target
[Unit]
Description=orora-cli.target for orora-cli quadlet

[Install]
WantedBy=default.target
